# ft/src/main/shell/crontab
#-----------------------------------------
## 每周2小时执行一次
* */2 * * * /ROOT/project/crawler/ft/crawler.sh start 2>&1 1>> /ROOT/project/crawler/ft/run.log